# crontab
* * * * * cd /nutrition-service && php artisan schedule:run &> /dev/null
